# helpers/systemd/terraria.service
[Unit]
Description=Terraria Game Server
Documentation=https://terraria.fandom.com/wiki/Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=atlas
Group=atlas
WorkingDirectory=/opt/terraria-server
ExecStart=/opt/terraria-server/terraria-start
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=terraria

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
